# cron/crontab
*/5 * * * * /usr/bin/python3 /app/generate_rules.py >> /var/log/generate_rules_cron.log 2>&1
